# deploy/systemd/xsp-killer-lane-a-monitor.service
[Unit]
Description=XSP Killer — Lane A morning exit monitor
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/xsp-killer
Environment=PYTHONPATH=/opt/xsp-killer
Environment=XSP_LANE_A_RH_POLL=false
# Phase 1: MCP reads + review-only exit dry-runs (incl. proof-of-life canary).
# LIVE_EXITS stays false — no orders placed until operator GO + funding + soak.
Environment=XSP_LANE_A_RH_MCP=true
Environment=XSP_LANE_A_LIVE_EXITS=false
Environment=XSP_LANE_A_PHASE1_CANARY=true
EnvironmentFile=-/opt/xsp-killer/.env
ExecStart=/opt/xsp-killer/scripts/lane_a_monitor_cron.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
