# deploy/systemd/xsp-killer-lane-b-monitor.timer
[Unit]
Description=XSP Killer Lane B — daily 08:00 ET

[Timer]
OnCalendar=*-*-* 08:00:00 America/New_York
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
