# deploy/systemd/xsp-killer-lane-a-intraday.timer
[Unit]
Description=XSP Killer Lane A — intraday every 15m 09:30–16:00 ET

[Timer]
OnCalendar=*-*-* 09:30:00 America/New_York
OnCalendar=*-*-* 09:45:00 America/New_York
OnCalendar=*-*-* 10:00:00 America/New_York
OnCalendar=*-*-* 10:15:00 America/New_York
OnCalendar=*-*-* 10:30:00 America/New_York
OnCalendar=*-*-* 10:45:00 America/New_York
OnCalendar=*-*-* 11:00:00 America/New_York
OnCalendar=*-*-* 11:15:00 America/New_York
OnCalendar=*-*-* 11:30:00 America/New_York
OnCalendar=*-*-* 11:45:00 America/New_York
OnCalendar=*-*-* 12:00:00 America/New_York
OnCalendar=*-*-* 12:15:00 America/New_York
OnCalendar=*-*-* 12:30:00 America/New_York
OnCalendar=*-*-* 12:45:00 America/New_York
OnCalendar=*-*-* 13:00:00 America/New_York
OnCalendar=*-*-* 13:15:00 America/New_York
OnCalendar=*-*-* 13:30:00 America/New_York
OnCalendar=*-*-* 13:45:00 America/New_York
OnCalendar=*-*-* 14:00:00 America/New_York
OnCalendar=*-*-* 14:15:00 America/New_York
OnCalendar=*-*-* 14:30:00 America/New_York
OnCalendar=*-*-* 14:45:00 America/New_York
OnCalendar=*-*-* 15:00:00 America/New_York
OnCalendar=*-*-* 15:15:00 America/New_York
OnCalendar=*-*-* 15:30:00 America/New_York
OnCalendar=*-*-* 15:45:00 America/New_York
OnCalendar=*-*-* 16:00:00 America/New_York
Persistent=true
RandomizedDelaySec=20

[Install]
WantedBy=timers.target
