# deploy/systemd/xsp-killer-lane-b-monitor.service
[Unit]
Description=XSP Killer — Lane B LEAPS hedge monitor
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/xsp-killer
Environment=PYTHONPATH=/opt/xsp-killer
Environment=XSP_LANE_A_RH_POLL=false
Environment=XSP_LANE_A_RH_MCP=false
Environment=XSP_LANE_A_LIVE_EXITS=false
EnvironmentFile=-/opt/xsp-killer/.env
ExecStart=/opt/xsp-killer/scripts/lane_b_monitor_cron.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
