# deploy/systemd/xsp-killer-lane-a-monitor.timer
[Unit]
Description=XSP Killer Lane A morning — 09:35/09:45/10:00/10:30 ET

[Timer]
OnCalendar=*-*-* 09:35:00 America/New_York
OnCalendar=*-*-* 09:45:00 America/New_York
OnCalendar=*-*-* 10:00:00 America/New_York
OnCalendar=*-*-* 10:30:00 America/New_York
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
