# deploy/systemd/xsp-killer-lane-a-monitor.service
[Unit]
Description=XSP Killer — Lane A morning exit monitor
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/xsp-killer
Environment=PYTHONPATH=/opt/xsp-killer
Environment=XSP_LANE_A_RH_POLL=false
# MCP reads + Phase 1 review dry-run (incl. proof-of-life canary).
# Auto-exit place path exists but is gated OFF: LIVE_EXITS=false means the
# monitor reviews exits but places no orders. Flip LIVE_EXITS=true only after
# funding + paper soak + operator GO. Emergency halt: XSP_LANE_A_KILL_SWITCH=true
# or `touch .local/KILL_SWITCH` (blocks all placement, cancels still allowed).
Environment=XSP_LANE_A_RH_MCP=true
Environment=XSP_LANE_A_LIVE_EXITS=false
Environment=XSP_LANE_A_PHASE1_CANARY=true
EnvironmentFile=-/opt/xsp-killer/.env
ExecStart=/opt/xsp-killer/scripts/lane_a_monitor_cron.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
