# deploy/systemd/xsp-killer-lane-a-monitor.timer
[Unit]
Description=XSP Killer Lane A premarket sell window — 08:00–09:30 ET (Mon-Fri)

[Timer]
OnCalendar=Mon..Fri *-*-* 08:00:00 America/New_York
OnCalendar=Mon..Fri *-*-* 08:15:00 America/New_York
OnCalendar=Mon..Fri *-*-* 08:30:00 America/New_York
OnCalendar=Mon..Fri *-*-* 08:45:00 America/New_York
OnCalendar=Mon..Fri *-*-* 09:00:00 America/New_York
OnCalendar=Mon..Fri *-*-* 09:15:00 America/New_York
OnCalendar=Mon..Fri *-*-* 09:30:00 America/New_York
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
