# deploy/systemd/xsp-killer-lane-a-monitor.timer
[Unit]
Description=XSP Killer Lane A exit poll — every 15m Mon–Fri (GTH/RTH/curb when open)

[Timer]
OnCalendar=Mon..Fri *-*-* *:00:00 America/New_York
OnCalendar=Mon..Fri *-*-* *:15:00 America/New_York
OnCalendar=Mon..Fri *-*-* *:30:00 America/New_York
OnCalendar=Mon..Fri *-*-* *:45:00 America/New_York
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
