# deploy/systemd/xsp-killer-lane-a-intraday.service
[Unit]
Description=XSP Killer — Lane A intraday BB/VWAP scan (paper)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/xsp-killer
Environment=PYTHONPATH=/opt/xsp-killer
Environment=XSP_LANE_A_RH_POLL=false
# See entry.service: live buy path is gated OFF (LIVE_ENTRIES=false); dry-run only.
Environment=XSP_LANE_A_RH_MCP=true
Environment=XSP_LANE_A_LIVE_ENTRIES=false
Environment=XSP_LANE_A_LIVE_EXITS=false
Environment=XSP_LANE_A_PAPER_ENTRY=true
EnvironmentFile=-/opt/xsp-killer/.env
ExecStart=/opt/xsp-killer/scripts/lane_a_intraday_cron.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
